# root/etc/services.d/lighttpd/run
#!/usr/bin/with-contenv bash

lighttpd -D -f /etc/lighttpd/lighttpd.conf
