# root/etc/s6-overlay/s6-rc.d/svc-ssh/run
#!/usr/bin/with-contenv bash
set -e

SSHD_ARGS="${SSHD_ARGS:- -De}";

exec \
    /usr/sbin/sshd \
    ${SSHD_ARGS} \
    -f /etc/ssh/sshd_config
