# root/etc/s6-overlay/s6-rc.d/svc-lighttpd/run
#!/usr/bin/with-contenv bash
set -e

LIGHTTPD_ARGS="${LIGHTTPD_ARGS:- -D}";

exec \
    lighttpd \
    ${LIGHTTPD_ARGS} \
    -f /etc/lighttpd/lighttpd.conf
